# punch-timer.service
[Unit]
Description=Punch Timer Web Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/punch_timer
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="SECRET_KEY=your-secret-key-here"
Environment="FLASK_DEBUG=False"
ExecStart=/usr/bin/python3 -c "from waitress import serve; from app import app, init_data_files; init_data_files(); serve(app, host='0.0.0.0', port=7777)"
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
